# deployment/fundpilot.service
[Unit]
Description=FundPilot FastAPI Web Application
After=network.target

[Service]
Type=simple
User=fundpilot
Group=fundpilot
WorkingDirectory=/opt/fundpilot/app
EnvironmentFile=/etc/fundpilot/fundpilot.env
Environment=HOME=/opt/fundpilot
Environment=PYTHONPATH=/opt/fundpilot/app
Environment=ENVIRONMENT=production
Environment=ENABLE_API_DOCS=false
ExecStart=/opt/fundpilot/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8004 --workers 1 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
